# ansible/roles/yubikey/templates/45-yubikey.rules
ACTION=="remove", ATTRS{idVendor}=="1050", ATTRS{idProduct}=="0407", RUN+="/usr/bin/loginctl lock-sessions"
